# acnexus/files/etc/uci-defaults/99-acnexus
#!/bin/sh
# First-run: generate config.json from UCI settings
[ -f /root/.ac_controller/config.json ] && exit 0

mkdir -p /root/.ac_controller/logs

uci -q get acnexus.@acnexus[0] > /dev/null 2>&1 || uci -q get acnexus.settings > /dev/null 2>&1 || exit 0

# 纯 Python 生成 config.json，避免 sed 特殊字符注入
python3 -c "
import json, os

def get(k, d=''):
    import subprocess
    r = subprocess.run(['uci', '-q', 'get', 'acnexus.@acnexus[0].' + k], capture_output=True, text=True)
    v = r.stdout.strip()
    if v: return v
    r = subprocess.run(['uci', '-q', 'get', 'acnexus.settings.' + k], capture_output=True, text=True)
    return r.stdout.strip() or d

cfg = {
    'api_key': get('api_key'),
    'qw_host': get('qw_host'),
    'baidu_key': get('baidu_key', ''),
    'weather_provider': get('weather_provider', 'qweather'),
    'weather_provider_set': True,
    'location': {
        'lat': float(get('location_lat', '39.9')),
        'lon': float(get('location_lon', '116.4')),
        'name': get('location_name', '\u672a\u8bbe\u7f6e'),
    },
    'typhoon_provider': get('typhoon_provider', 'nmc'),
    'typhoon_ac_off': get('typhoon_ac_off') == '1',
    'appearance_mode': 'system',
    'enabled': True,
    'devices': {'broadlink': {}, 'xiaomi_cloud': {}},
    'current_device_mac': '',
    'current_brand_type': 'broadlink',
}

os.makedirs(os.path.dirname('/root/.ac_controller/config.json'), exist_ok=True)
with open('/root/.ac_controller/config.json', 'w') as f:
    json.dump(cfg, f, indent=2, ensure_ascii=False)

print('AC-Nexus config initialized')
"
